# admin_scripts/fixit.crontab
# Minute Hour Day Month Weekday Command
30 1 * * * cd /home/fixit/data; source /home/fixit/fixit_env/bin/activate; python /home/fixit/bin/fixit_preprocess_script.py; scp -P 2222 /home/fixit/data/*.pkl zhan3249@142.1.200.146:/home/fixit/fixit-prod
30 4 * * * cd /home/fixit/data; source /home/fixit/fixit_env/bin/activate; python /home/fixit/bin/pcrs_inbound_data_loader.py
